# systemd/meeting-capture.service.example
[Unit]
Description=Union meeting audio capture
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/meeting-pipeline
ExecStart=/home/YOUR_USERNAME/meeting-pipeline/bin/record-meeting.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target