# systemd/meeting-postprocess.service
[Unit]
Description=Run WhisperX on latest meeting recording

[Service]
Type=oneshot
User=bojinda
WorkingDirectory=/home/bojinda/meeting-pipeline
ExecStart=/bin/bash -lc 'LATEST="$(cat /home/bojinda/meeting-pipeline/meeting-recordings/last_recording.txt)" && /home/bojinda/meeting-pipeline/bin/postprocess-meeting.sh "$LATEST"'
TimeoutStartSec=infinity