# systemd/meeting-postprocess.path.example
[Unit]
Description=Watch for new meeting recordings

[Path]
PathModified=/home/YOUR_USERNAME/meeting-pipeline/meeting-recordings/last_recording.txt

[Install]
WantedBy=multi-user.target