# systemd/watch-lesson-recordings.path.example
[Unit]
Description=Watch for new lesson recordings

[Path]
PathModified=/home/YOUR_USERNAME/meeting-pipeline/lesson-recordings/last_recording.txt

[Install]
WantedBy=multi-user.target